# .docker/cron/crontab
* * * * * /usr/local/bin/php /var/www/html/yii update-share-price/update-all
