# .docker/cron/crontab
30-59/2 7 * * 1-5 /usr/local/bin/php /var/www/html/yii update-share-price/update-all
0-59/2 8 * * 1-5 /usr/local/bin/php /var/www/html/yii update-share-price/update-all
0-20/2 9 * * 1-5 /usr/local/bin/php /var/www/html/yii update-share-price/update-all
45-59/2 10 * * 1-5 /usr/local/bin/php /var/www/html/yii update-share-price/update-all
0-59/2 11 * * 1-5 /usr/local/bin/php /var/www/html/yii update-share-price/update-all
0-45/2 12 * * 1-5 /usr/local/bin/php /var/www/html/yii update-share-price/update-all

@daily /usr/local/bin/php /var/www/html/yii update-share-price/update-all
